# supervisord.conf
[supervisord]
nodaemon=true

[program:shiny]
command=node /usr/local/bin/shiny-server
autorestart=true
