# supervisord.conf
[supervisord]
nodaemon=true

[program:sshd]
command=/usr/sbin/sshd -D
autorestart=true

[program:shiny]
command=node /usr/local/bin/shiny-server
autorestart=true
